# roles/contestant/backupwork.service
[Unit]
Description=Contestant work backup

[Service]
ExecStart=/usr/sbin/backupwork
